# automount-failure-notify@.service
[Unit]
Description=Notify on automount failure for %i

[Service]
Type=oneshot
ExecStart=%SCRIPT_DIR%/notify-mount-failure.sh %i
