# auto-unmount.service
[Unit]
Description=Auto-unmount webservers at evening
ConditionPathIsMountPoint=%MOUNT_POINT%

[Service]
Type=oneshot
# Kill any processes using the mount point
ExecStartPre=/bin/bash -c 'fuser -km %MOUNT_POINT% || true'
# Wait a moment for processes to terminate
ExecStartPre=/usr/bin/sleep 2
# Force unmount
ExecStart=/usr/bin/umount -f -l %MOUNT_POINT%
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
