# auto-unmount.timer
[Unit]
Description=Auto-unmount webservers at evening timer

[Timer]
OnCalendar=*-*-* 20:59:00
Persistent=false

[Install]
WantedBy=timers.target
